# docker/supervisord.conf
; supervisord runs the three Mosaic dev processes in one container.
;
; The ordering problem it has to solve is real: the Platform exits if it cannot
; reach PostgreSQL, and supervisord has no depends_on. `priority` only decides
; *start order*, not readiness — so the Platform waits for the database itself
; rather than relying on having been started second.
;
; Logs go to stdout/stderr so `docker logs` is the single place to look, which is
; most of the reason for collapsing three containers into one.

[supervisord]
; Stated explicitly rather than left to inheritance. The Go toolchain lives in
; /usr/local/go/bin and postgres' binaries in a versioned directory; a child
; process that loses either fails as "not found", which reads like a broken
; image rather than a broken PATH.
environment=PATH="/usr/local/go/bin:/usr/lib/postgresql/18/bin:/go/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; `supervisorctl restart platform` is the reason to keep this: restarting one
; process without taking the database down with it is the main thing a
; single-container layout would otherwise cost you.
[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:postgres]
command=/usr/lib/postgresql/%(ENV_PG_MAJOR)s/bin/postgres -D %(ENV_PGDATA)s
user=postgres
priority=10
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:platform]
; pg_isready rather than a sleep: the container is slower to start on a cold
; volume than a warm one, and a fixed delay is either too short (a failed boot)
; or too long (waiting for nothing) depending on which day it is.
command=bash -c "until pg_isready -h 127.0.0.1 -U mosaic -d mosaic -q; do sleep 1; done; exec go run ./cmd/mosaic-platform"
directory=/src/platform
priority=20
autorestart=true
stopasgroup=true
killasgroup=true
; The Platform shuts down gracefully — it drains the outbox before releasing its
; listeners — so a restart that does not wait races the old process for :8080 and
; the new one dies with "address already in use". Long enough to cover the drain,
; short enough not to hang a developer restarting in a loop.
stopwaitsecs=25
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:web]
; npm install runs on every start because node_modules is a volume that may be
; empty on first create; it is a no-op once populated. sdui-react is built
; before the Shell starts, since the Shell consumes its dist/ as a workspace
; peer rather than its source.
command=bash -c "npm install --workspaces --include-workspace-root --no-fund --no-audit && npm run build --workspace @mosaic-media/sdui-react && exec npm run dev --workspace mosaic-shell -- --host 0.0.0.0 --port 5173 --strictPort"
directory=/src/web
priority=30
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
